# old/legacy/etlrestart.service
[Unit]
Description=Restarts Enemy Territory Legacy server service

[Service]
ExecStart=/bin/systemctl restart etlserver.service

